# root/etc/s6-overlay/s6-rc.d/init-mod-infisical-cli/up
#!/bin/sh
# Empty file to satisfy S6-overlay for oneshot services
# No action needed as this is a oneshot service
exit 0